# systemd/log-time.timer
[Unit]
Description=Run log-time.sh every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
